# infra/files/caddy/caddy.service.d/override.conf
[Service]
EnvironmentFile=/etc/caddy.env
Restart=on-failure
RestartSec=2s
