# infra/files/infisical/infisical-agent.service
[Unit]
Description=Infisical Agent (secrets sync)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/infisical agent --config /etc/infisical-agent/agent.yaml
Restart=on-failure
RestartSec=10
MemoryMax=256M
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
